# crontab
*/15 * * * * /home/somik/resourceMonitor.sh &
*/5 * * * *  python3 /home/somik/getLoad.py &